# HomeWork_3/kostia_pkg/etc/systemd/system/my_script.service
[Unit]
Description=My Service

[Service]
ExecStart=usr/bin/my_script.sh
Restart=always

[Install]
WantedBy=multi-user.target
